# ansible/roles/caddy/templates/caddy.service.j2
[Unit]
Description=Caddy Reverse Proxy (HelmV2)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=exec
User={{ username }}
ExecStart=/usr/bin/caddy run --config {{ data_dir }}/caddy/Caddyfile
ExecReload=/usr/bin/caddy reload --config {{ data_dir }}/caddy/Caddyfile
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
